# ipv4-ipv6.ini
[program:ipv4-ipv6_uwsgi]
command = /usr/Python2.7/bin/uwsgi --init-paste /var/www/ipv4-ipv6.yosida95.com/production.ini
process_name = %(program_name)s
numproc = 2
autostart = true
user = www-data
stopsignal = INT
redirect_stderr = true
stdout_logfile = /var/www/ipv4-ipv6.yosida95.com/log/%(program_name)s.log
enviroment = PYTHON_EGG_CACHE=/var/www/ipv4-ipv6.yosida95.com/.cache
